# salt/prometheus/files/start/upstart/prometheus.conf
{% from "prometheus/prometheus.map" import prometheus with context %}

description "prometheus"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]
limit nofile 100000 100000

script
  [ -f /etc/default/$UPSTART_JOB ] && . /etc/default/$UPSTART_JOB
  
  exec {{ prometheus.dir }}/prometheus-{{ prometheus.version }}.linux-amd64/prometheus --config.file={{ prometheus.config_path }}/prometheus.yml --storage.tsdb.path={{ prometheus.data_dir }} --storage.tsdb.retention {{ prometheus.retention }}

end script

