# salt/alertmanager/files/start/upstart/alertmanager.conf
{% from "alertmanager/alertmanager.map" import alertmanager with context %}

description "Alertmanager"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]
limit nofile 100000 100000

script
  [ -f /etc/default/$UPSTART_JOB ] && . /etc/default/$UPSTART_JOB
  
  exec {{ alertmanager.dir }}/alertmanager-{{ alertmanager.version }}.linux-amd64/alertmanager -config.file={{ alertmanager.config_path }}/config.yml -mesh.peer={{ alertmanager.mesh_peer }}:6783 

end script

